# bash-config/profile.d/05-Python

#pathmunge PYTHONPATH "${HOME}/script/Python"

# Virtual env wrapper
if [[ -e /usr/local/bin/virtualenvwrapper.sh ]]; then
    . /usr/local/bin/virtualenvwrapper.sh
fi

function virtualenv-gt() {
    export VIRTUALENVWRAPPER_VIRTUALENV=virtualenv-gt
    export VIRTUALENVWRAPPER_PYTHON=python-gt
    . /usr/local/bin/virtualenvwrapper-gt.sh
}

# vim: filetype=sh:
