# bash-config/profile.d/05-Python

pathmunge PYTHONPATH "${HOME}/script/Python"

# Virtual env wrapper
if [[ -e /usr/local/bin/virtualenvwrapper.sh ]]; then
    . /usr/local/bin/virtualenvwrapper.sh
fi

function virtualenv-gt() {
    export VIRTUALENVWRAPPER_VIRTUALENV=virtualenv-gt
    export VIRTUALENVWRAPPER_PYTHON=python-gt
    . /usr/local/bin/virtualenvwrapper-gt.sh
}

if which rlwrap >/dev/null 2>&1; then
    function python() {
        if [[ $# -eq 0 ]]; then
            rlwrap -H ~/.python_history python
        else
            python "$@"
        fi
    }
    function python3() {
        if [[ $# -eq 0 ]]; then
            rlwrap -H ~/.python_history python3
        else
            python3 "$@"
        fi
    }
fi

# vim: filetype=sh:
