# bash-config/profile.d/45-Docker

# On macOS machines that cannot run 'native' Docker, a VM is run
# to host Docker. We need to source the env so that docker on
# macOS can transparently control the docker daemon running in
# the VM
if [[ -n "$OS_MAC" ]]; then

    if [[ -e /usr/local/bin/docker-machine ]]; then

        # All of the following operate on the 'default' VM

        #logterm "Sourcing default VM's docker config"
        eval "$(/usr/local/bin/docker-machine env >/dev/null 2>&1)"

        dockerssh() {
            /usr/local/bin/docker-machine ssh "$@"
        }

        dockerup() {
            /usr/local/bin/docker-machine start "$@"
        }

        dockerdown() {
            /usr/local/bin/docker-machine stop "$@"
        }
    fi
fi

# List mounted volumes for docker container
function docker-volumes() {
    if [[ $# -eq 0 ]]; then
        echo "Usage: docker-volumes <container>"
        return
    fi
    docker inspect --format '|Volumes:|--------|{{range .Mounts}}Host -> Container: {{.Source}} -> {{.Destination}}|{{end}}' $1 | tr '|' '\n' | tr -d '{' | tr -d '}'
}


# vim:filetype=sh:

